# .config/systemd/user/random_wallpaper.service
[Unit]
Description=Random Dank Material Shell's Wallpaper

[Service]
Type=oneshot
ExecStart="/home/nutty/.local/bin/random_wall"
